# backend/deploy/roam-symbolicator.service
[Unit]
Description=Roam MetricKit symbolicator worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=roam
Group=roam
WorkingDirectory=/var/lib/roam-symbolicator
EnvironmentFile=/etc/roam-symbolicator/env
ExecStart=/usr/local/bin/roam-backend --symbolicate
Restart=on-failure
RestartSec=30s

# Symbolication mmaps large dSYMs; raise the open-files cap.
LimitNOFILE=65536

# FUSE access. ipsw shells out to apfs-fuse to mount the cryptex DMG inside
# remote IPSWs so it can copy out the dyld_shared_cache.
DeviceAllow=/dev/fuse rw

# Hardening
# NoNewPrivileges is intentionally NOT set: libfuse execs the setuid
# fusermount3 helper to perform mount(2), which NoNewPrivileges would block.
# CapabilityBoundingSet caps what the setuid escalation can grant.
CapabilityBoundingSet=CAP_SYS_ADMIN
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
